# container/systemd/tuturu-coturn.service
[Unit]
Description=coturn TURN/STUN server for tuturu
Documentation=https://github.com/lobziik/tuturu
After=tuturu-certbot.service tuturu-init.service network.target
Requires=tuturu-certbot.service tuturu-init.service

[Service]
Type=simple
ExecStart=/usr/bin/turnserver -c /etc/coturn/turnserver.conf
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/coturn /run

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
