# container/systemd/tuturu-certbot-renew.service
[Unit]
Description=Renew Let's Encrypt certificates for tuturu
Documentation=https://github.com/lobziik/tuturu
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Stop nginx to free port 80, renew certs, then restart services
ExecStart=/usr/bin/certbot renew --quiet \
    --pre-hook "systemctl stop tuturu-nginx" \
    --post-hook "systemctl start tuturu-nginx && systemctl reload tuturu-coturn"

StandardOutput=tty
StandardError=tty
TTYPath=/dev/console
TTYVTDisallocate=yes
