# container/systemd/tuturu-app.service
[Unit]
Description=tuturu signaling server (Bun WebSocket)
Documentation=https://github.com/lobziik/tuturu
After=tuturu-init.service tuturu-redis.service network.target
Requires=tuturu-init.service
# Soft dependency on Redis - works without it (credentials just expire naturally)
Wants=tuturu-redis.service

[Service]
Type=simple
EnvironmentFile=/etc/tuturu/env
ExecStart=/usr/local/bin/tuturu-server
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/tuturu

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/tuturu

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
