# container/systemd/tuturu-certbot-renew.timer
[Unit]
Description=Daily certificate renewal check for tuturu
Documentation=https://github.com/lobziik/tuturu

[Timer]
# Run daily at a random time within 1 hour window
OnCalendar=daily
RandomizedDelaySec=3600
# Catch up on missed runs
Persistent=true

[Install]
WantedBy=timers.target
