# container/systemd/tuturu-redis.service
[Unit]
Description=Redis for tuturu TURN credential revocation
After=network.target
Before=tuturu-coturn.service tuturu-app.service

[Service]
Type=simple
# Redis in-memory only, no persistence needed (blacklist entries have TTL)
# Bind to localhost only (container-internal communication)
# 64MB max memory with LRU eviction for safety
ExecStart=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --daemonize no --maxmemory 64mb --maxmemory-policy allkeys-lru --save ""
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
